# script/cron.d/block_hash.cron
MAILTO=ibn1029@gmail.com
*/15 * * * * viage /bin/bash /Users/viage/Work/App/block_hash/script/do.sh 
